# infra-images/sandbox/services.d/dockerd/run
#!/bin/sh
exec dockerd --host unix:///var/run/docker.sock --log-format=json 2>&1
