# infra-images/sandbox/services.d/dns-proxy/run
#!/bin/sh
if [ -f /var/lib/dboxed/sandbox.env ]; then
  source /var/lib/dboxed/sandbox.env
fi
exec dboxed sandbox dns-proxy
