# infra-image/services.d/dboxed/run
#!/bin/sh
if [ -f /var/lib/dboxed/sandbox.env ]; then
  source /var/lib/dboxed/sandbox.env
fi
exec dboxed box run-in-sandbox 2>&1
