# infra-image/services.d/dboxed/run
#!/bin/sh
exec dboxed run-box-in-sandbox 2>&1
