# infra-image/services.d/dboxed/run
#!/bin/sh
exec dboxed box run-in-sandbox 2>&1
